# projects/14-observability/terraform/templates/observability-app.service
[Unit]
Description=Observability lab web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/observability-app/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target